# docker/packaging/prerm
#!/bin/bash
set -e
systemctl stop buddy.service || true
systemctl disable buddy.service || true
